# cluster_images/manager/supervisord.conf
[supervisord]
nodaemon=true

[program:bro]
command=/bro/bin/run_bro.sh

[program:sshd]
command=/usr/sbin/sshd -D -p 2022
